# experiments/01-measuring-memory-consumption-of-python-programs/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
username=supervisord
password=supervisord

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
username=supervisord
password=supervisord

[program:monitor_memory_usage]
command=/app/scripts/monitor_memory_usage.sh
autostart=true
autorestart=false
startretries=0
startsecs=1
stdout_logfile=%(ENV_MEMORY_USAGE_LOG_FILE)s
stderr_logfile=%(ENV_MEMORY_USAGE_LOG_FILE)s
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:monitor_page_faults]
command=/app/scripts/monitor_page_faults.sh
autostart=true
autorestart=false
startretries=0
startsecs=1
stdout_logfile=%(ENV_PAGE_FAULTS_LOG_FILE)s
stderr_logfile=%(ENV_PAGE_FAULTS_LOG_FILE)s
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:main]
command=/bin/sh -c "python3 /app/experiment/main.py ${RUN_COMMAND}"
autostart=true
autorestart=false
startretries=0
startsecs=1
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[eventlistener:processes]
command=bash -c "echo READY && read line && kill -SIGQUIT $PPID"
events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL