# experiments/.bak/01-measuring-memory-consumption-of-python-programs/supervisord.conf
[supervisord]
nodaemon=true
logfile=%(ENV_SUPERVISORD_LOG_FILE)s
pidfile=%(ENV_SUPERVISORD_PID_FILE)s
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
minfds=1024
minprocs=200
user=root

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
username=supervisord
password=supervisord

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
username=supervisord
password=supervisord

[program:monitor_memory_usage]
command=/app/scripts/monitor_memory_usage.sh
autostart=true
autorestart=false
startretries=0
startsecs=0
exitcodes=0
stdout_logfile=%(ENV_MEMORY_USAGE_LOG_FILE)s
stderr_logfile=%(ENV_MEMORY_USAGE_LOG_FILE)s
priority=1

[program:monitor_page_faults]
command=/app/scripts/monitor_page_faults.sh
autostart=true
autorestart=false
startretries=0
startsecs=0
exitcodes=0
stdout_logfile=%(ENV_PAGE_FAULTS_LOG_FILE)s
stderr_logfile=%(ENV_PAGE_FAULTS_LOG_FILE)s
priority=2

[program:memory_pressure_monitor]
command=/app/scripts/monitor_memory_pressure.sh
autostart=true
autorestart=false
startretries=0
startsecs=0
exitcodes=0
stdout_logfile=%(ENV_MEMORY_PRESSURE_LOG_FILE)s
stderr_logfile=%(ENV_MEMORY_PRESSURE_LOG_FILE)s
priority=3

[program:limit_memory_usage]
command=/app/scripts/limit_memory_usage.sh
autostart=true
autorestart=false
startretries=0
startsecs=0
exitcodes=0
stdout_logfile=%(ENV_LIMIT_MEMORY_USAGE_LOG_FILE)s
stderr_logfile=%(ENV_LIMIT_MEMORY_USAGE_LOG_FILE)s
priority=4

[program:main]
command=/bin/sh -c "python3 /app/experiment/main.py ${RUN_COMMAND} || echo $? > /tmp/main_exit_code"
autostart=true
autorestart=false
startretries=0
startsecs=0
exitcodes=0
stopsignal=TERM
stdout_logfile=%(ENV_MAIN_LOG_FILE)s
stderr_logfile=%(ENV_MAIN_LOG_FILE)s
priority=5


[eventlistener:quit_on_main_exit]
command=/bin/bash /app/scripts/quit_on_main_exit.sh
events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL